# supervisor/conf.d/gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/A2/BookingPortal
command=/home/ubuntu/env/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/A2/BookingPortal/app.sock BookingPortal.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log

[group:guni]
programs:gunicorn
